# week7/only_bg/start_at.crontab
30 12 * * * env DISPLAY=:0 /home/rado_sz/Programming101-3/week7/only_bg/task.sh